# files/etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
pidfile=/run/supervisord/supervisord.pid
logfile=/var/log/supervisord.log

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:pwm]
command=java -cp /usr/share/pwm/libs/* -jar /usr/share/pwm/pwm-onejar-1.8.0-SNAPSHOT.jar -applicationPath /config
autorestart=true
